# loadbalance/config/supervisor/xc_vm.conf
[program:xc_vm]
command=/scripts/init-xc_vm.sh
autostart=true
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
